# devops/deb-config/vdi-backend/root/DEBIAN/postinst
#!/bin/

set -e

config_baloon()
{
	local lsmod_exit

	lsmod_exit=$(lsmod | grep -i balloon)

	if [ -z "$lsmod_exit" ]; then
		LOG "* ${CYAN}${MY_APPLICATION}: ${WHITE}balloon module has been loaded${NORMAL}"
		modprobe virtio_balloon
	fi
	
	# добавляем модуль в /etc/modules-load.d/modules.conf
	if ! grep -q "balloon" /etc/modules-load.d/modules.conf; then
		echo "virtio_balloon" >> /etc/modules-load.d/modules.conf
	fi

}

configure_redis()
{
	systemctl enable redis-server.service
	echo 'requirepass 4NZ7GpHn4IlshPhb' | sudo tee -a /etc/redis/redis.conf
	sed -i 's/bind 127.0.0.1/bind 0.0.0.0/g' /etc/redis/redis.conf
	systemctl restart redis-server
}

case "$1" in
	configure)
		# действия и при установке и при обновлении

		config_zfs
		config_dlm
		config_baloon
		systemctl daemon-reload || true

		# действия только при установке
		if [ ! -n "$2" ] ; then
			configure_redis
		fi

		exit 0
		;;
	abort-upgrade|abort-remove|abort-deconfigure)
		systemctl restart cron
		exit 0
		;;
	*)
		LOG "postinst called with unknown argument: $*" >&2
		exit 1
		;;
esac
